# deployment/systemd/honeypot-airdrop.service
[Unit]
Description=Honeypot AirDrop Threat Detection System
Documentation=man:systemd.service(5)
After=network.target docker.service
Requires=network.target

[Service]
Type=simple
Restart=always
RestartSec=10
User=root
WorkingDirectory=/home/burner/honeypot-minimal
ExecStart=/usr/bin/python3 /home/burner/honeypot-minimal/airdrop_threat_detector.py
StandardOutput=append:/home/burner/honeypot-minimal/logs/airdrop_threats.log
StandardError=append:/home/burner/honeypot-minimal/logs/airdrop_threats.log

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=ELASTICSEARCH_HOST=localhost:9200
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=honeypot.target