# deployment/systemd/honeypot-wifi-detector.service
[Unit]
Description=Honeypot WiFi Threat Detector
After=network.target docker.service
Requires=network.target

[Service]
Type=simple
Restart=always
RestartSec=10
User=root
WorkingDirectory=/home/burner/honeypot-minimal
ExecStart=/usr/bin/python3 /home/burner/honeypot-minimal/wifi_threat_detector_filtered.py
StandardOutput=append:/home/burner/honeypot-minimal/logs/wifi_threat_detector.log
StandardError=append:/home/burner/honeypot-minimal/logs/wifi_threat_detector.log

[Install]
WantedBy=honeypot.target