# honeyman-v2/deployment/ops/honeyman-healthcheck.service
[Unit]
Description=Honeyman uptime probe (one shot)
Documentation=https://github.com/SecurityWard/honeyman-Project

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/honeyman-healthcheck.sh
# Operators override API_BASE / TIMEOUT / HEALTHCHECK_WEBHOOK via
# /etc/honeyman/healthcheck.env. Leading `-` means "skip silently
# if the file isn't there" so the unit still starts on hosts that
# haven't configured a webhook yet.
EnvironmentFile=-/etc/honeyman/healthcheck.env
